# roles/www_ddmail_se/templates/ddmail_gunicorn.service.j2
[Unit]
Description=Gunicorn instance to serve ddmail
After=network.target

[Service]
User=www-ddmail
Group=www-ddmail
WorkingDirectory=/opt/ddmail_venv/
Environment="PATH=/opt/ddmail_venv/bin/"
ExecStart=/opt/ddmail_venv/bin/gunicorn --workers 2 "ddmail:create_app()"

[Install]
WantedBy=multi-user.target
